# rootfs/etc/init.d/pulseaudio
#!/bin/sh
### BEGIN INIT INFO
# Provides:          pulseaudio
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Enable service
# Description:       Enable service by init script.
### END INIT INFO

PATH="/bin:/usr/bin:/sbin:/usr/sbin"
DESC="PulseAudio sound server"
NAME="pulseaudio"
DAEMON="/usr/bin/pulseaudio"
DAEMONOPTS="--start --disallow-exit --exit-idle-time=-1"
SCRIPTNAME="/etc/init.d/$NAME"
EXEUSER="root"

if [ -f /etc/default/webos ]; then
    set -a
    . /etc/default/webos
    set +a
fi

. /lib/lsb/init-functions

case "$1" in
    start)
        log_daemon_msg "Starting service" "$NAME"
        # Check if already running based on process query
        if pgrep -u $EXEUSER -xf "$DAEMON $DAEMONOPTS" > /dev/null 2>&1; then
            log_end_msg 0
            exit 0
        fi
        # Start the daemon
        start-stop-daemon --start --background --user $EXEUSER --chuid $EXEUSER --exec $DAEMON -- $DAEMONOPTS
        # Wait a moment for the process to start
        sleep 1
        # Verify the process started successfully
        if pgrep -u $EXEUSER -xf "$DAEMON $DAEMONOPTS" > /dev/null 2>&1; then
            log_end_msg 0
        else
            log_end_msg 1
        fi
        ;;
    stop)
        log_daemon_msg "Stopping service" "$NAME"
        # Kill the process based on process query
        if pgrep -u $EXEUSER -xf "$DAEMON $DAEMONOPTS" > /dev/null 2>&1; then
            pkill -u $EXEUSER -xf "$DAEMON $DAEMONOPTS" 2>/dev/null
            # Wait for the process to terminate
            sleep 1
            # Force kill if still running
            if pgrep -u $EXEUSER -xf "$DAEMON $DAEMONOPTS" > /dev/null 2>&1; then
                pkill -9 -u $EXEUSER -xf "$DAEMON $DAEMONOPTS" 2>/dev/null || true
            fi
        fi
        log_end_msg 0
        ;;
    restart)
        log_daemon_msg "Restarting service" "$NAME"
        $0 stop
        sleep 1
        $0 start
        ;;
    status)
        # Check status based on process query only
        PID=$(pgrep -u $EXEUSER -xf "$DAEMON $DAEMONOPTS" 2>/dev/null)
        if [ -n "$PID" ]; then
            log_daemon_msg "$NAME is running (PID: $PID)"
            log_end_msg 0
            exit 0
        else
            log_daemon_msg "$NAME is not running"
            log_end_msg 3
            exit 3
        fi
        ;;
    *)
        echo "Usage: $0 {start|stop|status|restart}"
        exit 1
        ;;
esac
exit 0
