# rootfs/etc/s6-overlay/s6-rc.d/vncserver/run
#!/command/execlineb -P
s6-envdir /run/s6/container_environment
s6-envdir /etc/s6-overlay/s6-rc.d/vncserver/env
foreground { 
    sh -c "rm -rf /tmp/.X*-lock /tmp/.X11-unix" 
}
importas -SD "root" USER
s6-setuidgid ${USER}
foreground {
    sh -c "
        if [ ! -f \"\$HOME/.kasmpasswd\" ]; then
            printf \"%s\\n%s\\n\" \"\$NOVNC_PASSWORD\" \"\$NOVNC_PASSWORD\" | vncpasswd -u \"\$NOVNC_USER\" -o -w -r
        fi
    "
}
sh -c "/usr/bin/vncserver -fg -select-de XFCE -SecurityTypes None -UseIPv6 0 -sslOnly 0 -localhost 0 -websocketPort 6901 \${NOVNC_IP:+-publicIP \"\$NOVNC_IP\"} :1"
